# com.mole.backgroundscan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mole.backgroundscan</string>
    <key>ProgramArguments</key>
    <array>
        <!-- Binary name matches build_release_dmg.sh: binary is copied as "Mole" not "MoleApp" -->
        <string>/Applications/Mole.app/Contents/MacOS/Mole</string>
        <string>--background-scan</string>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>StartInterval</key>
    <!-- 1800s = 30 minutes default. BackgroundScheduleManager.swift writes a modified copy at runtime for different intervals. -->
    <integer>1800</integer>
    <key>StandardOutPath</key>
    <!-- Note: Full path expansion ($HOME) is not supported in plists; installer must write the actual path. -->
    <string>/tmp/com.mole.backgroundscan.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/com.mole.backgroundscan.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>MOLE_BACKGROUND_MODE</key>
        <string>1</string>
    </dict>
</dict>
</plist>
